# .crontab
0 0 * * * npm run import:textsAuthors
0 1 * * * npm run import:texts
